# deploy/supervisor/feishu-bot.conf
# Supervisor configuration for Feishu Bot
#
# 安装步骤：
# 1. 复制此文件到 /etc/supervisor/conf.d/feishu-bot.conf
# 2. 修改 command、directory、user 等配置
# 3. 运行：sudo supervisorctl reread
# 4. 运行：sudo supervisorctl update
# 5. 启动：sudo supervisorctl start feishu-bot
#
# 管理命令：
# - 查看状态：sudo supervisorctl status feishu-bot
# - 查看日志：sudo supervisorctl tail -f feishu-bot
# - 重启服务：sudo supervisorctl restart feishu-bot
# - 停止服务：sudo supervisorctl stop feishu-bot

[program:feishu-bot]

# 启动命令（修改为实际路径）
# 如果使用虚拟环境：/path/to/venv/bin/python
# 如果使用系统 Python：/usr/bin/python3
command=/usr/bin/python3 -m src.main

# 工作目录（修改为实际路径）
directory=/path/to/cc-feishu

# 运行用户（修改为实际用户）
user=your-user

# 自动启动和重启
autostart=true
autorestart=true

# 启动优先级（数字越小越先启动）
priority=999

# 进程启动后等待时间（秒）
startsecs=10

# 启动失败重试次数
startretries=3

# 停止信号
stopsignal=TERM

# 停止等待时间（秒）
stopwaitsecs=10

# 重定向 stderr 到 stdout
redirect_stderr=true

# 日志配置
stdout_logfile=/var/log/feishu-bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10

# 环境变量
# 注意：敏感信息应该放在 .env 文件中，不要写在这里
environment=PATH="/usr/local/bin:/usr/bin:/bin",PYTHONUNBUFFERED="1"

# 子进程管理
# WebSocket 断开时进程会退出，supervisor 会自动重启
killasgroup=true
stopasgroup=true
